# mail-poller/mail-poller.service
[Unit]
Description=Get all new mail

[Service]
Type=simple
TimeoutStartSec=110
User=frame
Group=frame
ExecStart=/usr/bin/offlineimap
